# mkosi.images/base/mkosi.extra/usr/lib/systemd/system/systemd-sysupdate.service.d/statedir.conf
[Service]
StateDirectory=extensions extensions.available
