# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/sysroot.conf
[Unit]
After=sysroot.mount
ConditionDirectoryNotEmpty=|/sysroot/usr/lib/repart.d
ConditionKernelCommandLine=!mangos.install=destructive

[Service]
# Ensure systemd-repart has access to the machine-id so that it
# will generate correct partition UUID's for /var
ExecStartPre=mount -o bind /etc/machine-id /sysroot/etc/machine-id
ExecStartPost=umount /sysroot/etc/machine-id
