# mkosi.images/initrd/mkosi.extra/usr/lib/systemd/system/systemd-repart.service.d/tpm.conf
[Unit]
After=systemd-tpm2-setup-early.service sysroot.mount

[Service]
ExecStart=

# PCR 7:  Secure Boot state; changes when UEFI SecureBoot mode is enabled/disabled,
#         or firmware certificates (PK, KEK, db, dbx, ...) changes.
# PCR 11: systemd-stub(7) measures the ELF kernel image, embedded initrd and other
#         payload of the PE image it is placed in into this PCR. systemd-
#         pcrphase.service(8) measures boot phase strings into this PCR at various 
#         milestones of the boot process.
#
# Binding the key to PCR 11 means it'll only be usable if the measured values of
# the loaded UKI match the embedded signed policy, signed by our key.
#
# Binding the key to PCR 7 means it will NOT be available if the Secure Boot
# policy has been modified.
ExecStart=systemd-repart --dry-run=no --tpm2-pcrs=7 --tpm2-public-key-pcrs=11 --tpm2-pcrlock=
