# mkosi.extra/usr/lib/systemd/system/systemd-networkd.service.d/after-boot.conf
[Unit]
WantsMountsFor=/boot/default.network
