# mkosi.extra/usr/lib/systemd/system/mangos-sd-gh-proxy.service
[Unit]
Description=Let sysupdate pull stuff from Github

[Service]
ExecStart=/usr/bin/mangos-sd-gh-proxy
