# mkosi.extra/usr/lib/systemd/system/mangos-save-boot.service
[Unit]
Description=Save boot choice
After=boot-complete.target

[Service]
Type=oneshot
ExecStart=/usr/sbin/grub-set-default %M_%A
