# mkosi.extra/usr/lib/systemd/system/ssh.service.d/ssh-keygen.conf
[Unit]
Wants=sshd-keygen.target
After=sshd-keygen.target
